# resources/aprs_wx/aprs_wx.service.in
[Unit]
Description=APRS_WX weather aggregation execution service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=${DAEMON_USER}
Group=${DAEMON_GROUP}
ExecStart=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_BINDIR}/${PROJECT_NAME}

[Install]
WantedBy=multi-user.target
