# cmake-build/_CPack_Packages/Linux/DEB/APRS_WX-0.3.2-Linux/APRS_WX_SYSTEMD/postinst
#!/bin/bash
systemctl daemon-reload
chown daemon:daemon /usr/local/share/APRS_WX/config.txt
if systemctl is-enabled aprs_wx.service > /dev/null; then
    systemctl restart aprs_wx.service
fi
